# wpclearcache.service
[Service]
ExecStart=/usr/local/bin/wpclearcache

[Install]
WantedBy=default.target